# contrib/new-boot-0.2/etc/inittab
#NOTE: currently init has hardcoded a limit of 30 tty lines and a limit
#of 150 characters for every line (tty and other).

bootprog     =  rc 3
fileprefix   =  /sbin/init.d/
PATH         =  /sbin:/bin:/usr/sbin:/usr/bin
INIT_PATH    =  .:runlevel.1:runlevel.2:runlevel.3:runlevel.4:runlevel.5

tty1:linux:/sbin/agetty tty1 9600
tty2:linux:/sbin/agetty tty2 9600
tty3:linux:/sbin/agetty tty3 9600
tty4:linux:/sbin/agetty tty4 9600
tty5:linux:/sbin/agetty tty5 9600
tty6:linux:/sbin/agetty tty6 9600
                                         		         
finalprog    =
quietrollback=yes
ctrlaltdel   =/bin/sh -c "PATH=$INIT_ROOT:$PATH; echo $'\\a\\nAuf Wiedersehen!\\n' ; initctl -u -r -f"
#ctrlaltdel   =/sbin/shutdown -r -q now
