# udev-config/27-firmware.rules
# This rule simply calls the firmware_helper program to load firmware

ACTION=="add", SUBSYSTEM=="firmware", RUN+="/lib/udev/firmware_helper"
