# udev-config/95-udev-late.rules
# Send all events to this Unix-domain socket, to be caught by udevmonitor
RUN+="socket:/org/kernel/udev/monitor"

